# cleanupVolumes.cron
* * * * * /usr/local/bin/node /app/dist/cleanupVolumes.js >> /var/log/cron.log 2>&1
